# .config/autostart/unclutter-hide.desktop
[Desktop Entry]
Type=Application
Name=unclutter-hide
Exec=/usr/bin/unclutter -idle 2

