# .config/autostart/xbacklight-default.desktop
[Desktop Entry]
Type=Application
Name=xbacklight-default
Exec=/usr/bin/xbacklight -set 4

